# ml-waf/ml-waf.service
[Unit]
Description=ML-Enhanced WAF FastAPI Daemon
After=network.target redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/ml-waf
ExecStart=/opt/ml-waf/venv/bin/uvicorn ml_server:app --host 127.0.0.1 --port 8003 --workers 2
Restart=always
RestartSec=5

# Process resource limitations
# MemoryMax is set to 1G to safely support loading 2 worker process model binaries in memory (numpy, xgboost, scikit-learn)
MemoryMax=1G
MemoryHigh=900M
CPUQuota=100%

# Security configurations
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
